# tools/searxng/searxng.service
[Unit]
Description=SearXNG Meta Search Engine
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/podman rm -f searxng
ExecStart=/usr/bin/podman run --name searxng \
  -p 127.0.0.1:47123:8080 \
  -v %h/.config/searxng/settings.yml:/etc/searxng/settings.yml:ro \
  -e SEARXNG_BASE_URL=http://localhost:47123/ \
  --cap-drop ALL \
  docker.io/searxng/searxng:latest
ExecStop=/usr/bin/podman stop searxng
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
