# ai/xbol-observability/quadlet/xbol-grafana.container
[Unit]
Description=XBOL Grafana dashboard
After=xbol-victoriametrics.service xbol-loki.service
Requires=xbol-victoriametrics.service xbol-loki.service

[Container]
ContainerName=xbol-grafana
Image=docker.io/grafana/grafana-oss:latest
Pull=never
Network=xbol-observability.network
PublishPort=127.0.0.1:43000:3000
Volume=xbol-grafana-data:/var/lib/grafana
Volume=%h/.dotfiles/ai/xbol-observability/grafana/provisioning:/etc/grafana/provisioning:ro
Volume=%h/.dotfiles/ai/xbol-observability/grafana/dashboards:/var/lib/grafana/dashboards:ro
Environment=GF_SERVER_HTTP_ADDR=0.0.0.0
Environment=GF_SECURITY_ADMIN_USER=admin
Environment=GF_SECURITY_ADMIN_PASSWORD=xbol-local-admin-change-me
Environment=GF_AUTH_ANONYMOUS_ENABLED=true
Environment=GF_AUTH_ANONYMOUS_ORG_ROLE=Viewer
Environment=GF_USERS_ALLOW_SIGN_UP=false
Environment=GF_ANALYTICS_REPORTING_ENABLED=false
Environment=GF_ANALYTICS_CHECK_FOR_UPDATES=false
User=0
DropCapability=all
NoNewPrivileges=true

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
