# ai/openchamber/openchamber-serve.service
[Unit]
Description=OpenChamber Web Server
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=uwsm app -- /home/carlos/.local/bin/openchamber-serve.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
