# packages/podman/systemd/podman-system-prune.timer
[Unit]
Description=Weekly rootless Podman system prune

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
