# ai/opencode/opencode-serve.service
[Unit]
Description=OpenCode Web Server
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=uwsm app -- /home/carlos/.opencode/bin/opencode serve --port 3014 --hostname 127.0.0.1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
