# ai/opencode/opencode-serve.service
[Unit]
Description=OpenCode Web Server
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=uwsm app -- /home/carlos/.local/bin/opencode-serve.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
