# ai/xbol-observability/systemd/xbol-gcp-log-fetch.service
[Unit]
Description=Fetch narrowly scoped XBOL GCP logs into local spool files
Documentation=file:%h/.dotfiles/ai/xbol-observability/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/xbol-gcp-log-fetch --once
