# tools/hyprspace/hyprspace-watch.service
[Unit]
Description=hyprspace stray-window watcher
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=uwsm app -- %h/.local/bin/hyprspace watch
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
