# languages/python/uv/systemd/uv-cache-prune.service
[Unit]
Description=Prune unreachable uv cache entries

[Service]
Type=oneshot
ExecStart=/usr/bin/uv cache prune
IOSchedulingClass=idle
