# ai/openviking/quadlet/openviking.container
[Unit]
Description=OpenViking local context database

[Quadlet]
DefaultDependencies=false

[Container]
ContainerName=openviking
Image=ghcr.io/volcengine/openviking:latest
Pull=never
Network=host
Environment=OPENVIKING_CONFIG_FILE=/app/.openviking/ov.conf
Environment=OPENVIKING_SERVER_HOST=127.0.0.1
Environment=OPENVIKING_WITH_BOT=0
Volume=%h/.openviking:/app/.openviking:Z
Volume=%h/.local/share/openviking:/app/.openviking/data:Z
NoNewPrivileges=true
DropCapability=all

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
